# etc/sv/server/finish
#!/bin/sh

sleep ${SSH_SERVER_RESTART_DELAY}
