# etc/sv/client/finish
#!/bin/sh

sleep ${SSH_CLIENT_RESTART_DELAY}
