# etc/sv/client/run
#!/bin/sh

echo running command:
echo ${SSH_CLIENT_COMMAND}
echo
su -s /bin/sh ${SSH_CLIENT_USER} -c "${SSH_CLIENT_COMMAND}"
