# etc/sv/client/run
#!/bin/sh

echo running command:
echo ${SSH_CLIENT_COMMAND}
echo
${SSH_CLIENT_COMMAND}
